# supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:alfresco]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
command=/alfresco/init.sh

[program:nginx]
command=nginx
